# xxxsupervisord-kong.conf
[supervisord]
nodaemon=true

[program:pm2]
command=pm2 start /home/sttrweb/Oncoscape/server/app.js -i 0 --no-daemon

[program:opencpu]
command=service opencpu restart && tail -F /var/log/opencpu/apache_access.log

[program:kong]
command=kong start -c /etc/kong/kong.conf --nginx-conf /usr/local/kong/nginx-kong-oncoscape.template
priority=900
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true